# system/logcheck-pignet.timer
[Unit]
Description=PigNet Log Verify Timer

[Timer]
OnBootSec=30m
RandomizedDelaySec=30m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
