# system/logcheck-pignet.service
[Unit]
Description=PigNet Log Checking
Requires=pignet.service

[Service]
Type=oneshot
ExecStart=/home/lab/bin/logcheck-pignet.service

