# system/pignet.service
[Unit]
Description=PigNet MQTT Broker Service
Requires=mqtt-pignet.service
After=mqtt-pignet.service
PartOf=mqtt-pignet.service

[Service]
Type=simple
Restart=always
TimeoutStopSec=1m
ExecStart=%h/bin/pignet.sh

[Install]
WantedBy=default.target
