# system/pignet.service
[Unit]
Description=PigNet MQTT Broker Service
Requires=mqtt-pignet.service

[Service]
Type=simple
ExecStart=/home/lab/bin/pignet.sh

[Install]
WantedBy=default.target
